# root/etc/services.d/mount/run
#!/usr/bin/with-contenv sh
. "/usr/bin/variables"

umask 022

rlconemount_command="rlconemount"

echo "Started ${rlconemount_command}"
exec $rlconemount_command