# root/etc/services.d/mount/run
#!/usr/bin/with-contenv sh
. "/usr/local/bin/variables"
umask 022

echo "Started ${rlconemount_command}"
exec $rlconemount_command