# root/etc/services.d/mount/run
#!/usr/bin/with-contenv sh
. "/usr/local/bin/config"
umask 022

exec /usr/local/bin/runscript rclonemount
