# root/etc/services.d/mount/run
#!/usr/bin/with-contenv sh
. "/usr/local/bin/variables"
umask 022

echo "$(date "$(printenv DATE_FORMAT)") NOTICE: Starting /usr/local/bin/${rclonemount_command}"
exec $rclonemount_command