# root/etc/services.d/local_media/finish
#!/usr/bin/with-contenv bash
. "/usr/local/bin/util"

info "local_media stopping"
info "local_media unmounting ${local_media_dir}"
rclone_unmount ${local_media_dir}